# dist/systemd/zrepl_timers.target
[Unit]
Description=Zrepl Timers
DefaultDependencies=no
